# docker/supervisor/redis.ini
[program:redis]
command=redis-server
directory=/var/redis
